# systemd/eink-poster.service
[Unit]
Description=Motivational Poster

[Service]
User=ben
WorkingDirectory=/home/ben/eink-planner
ExecStart=python3 eink_poster.py

[Install]
WantedBy=multi-user.target