# systemd/motivational-poster.timer
[Timer]
Unit=motivational-poster
OnCalendar=*-*-* 17:00:00
Persistent=true