# systemd/motivational-poster.timer
[Timer]
Unit=motivational-poster.service
OnCalendar=*-*-* 17:15:00

[Install]
WantedBy=timers.target